# projects/09-24-7-ops/systemd/docker-compose@.service
[Unit]
Description=Docker Compose project at %I
Requires=docker.service network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
# %I is the systemd-escaped absolute path of the Compose directory.
WorkingDirectory=%I
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=300
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
