# projects/09-24-7-ops/systemd/pi-watchdog.service
[Unit]
Description=Pi watchdog — reboots if docker becomes unresponsive
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pi-watchdog.sh
